# ewatv-playout-backend/deployments/podman/quadlets/ewatv-redis.container
[Unit]
Description=EWATV Redis (playout cache)
After=network-online.target

[Container]
ContainerName=ewatv-redis
Image=docker.io/library/redis:7-alpine
Network=ewatv-net.network
Exec=redis-server --save "" --appendonly no --maxmemory 256mb --maxmemory-policy allkeys-lru
HealthCmd=redis-cli ping
AutoUpdate=registry

[Service]
Restart=always

[Install]
WantedBy=default.target
